# system/ipv6.service
[Unit]
Description=Discard unused blocks

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ipv6.sh
